# roles/ansible-pull/files/ansible-pull.timer
[Unit]
Description=Run ansible-pull everyday and on boot

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true


[Install]
WantedBy=timers.target
