# roles/poweroff-schedule/files/poweroff-schedule.service
[Unit]
Description = Poweroff at timer

[Service]
Type=oneshot
ExecStart=/usr/bin/poweroff
